# scripts/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Rectella Shopify webhooks (named tunnel: rectella-shopify)
After=rectella.service
Wants=rectella.service

[Service]
Type=simple
# Named tunnel — uses ~/.cloudflared/config.yml (hostname rectella.ctrlaltinsight.co.uk).
# DNS is fixed; no per-restart webhook re-registration needed.
ExecStart=/usr/bin/cloudflared tunnel run rectella-shopify
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
